# examples/hermes-vipd.service
[Unit]
Description=Hermes VIP Daemon
After=network.target

[Service]
Type=simple
User=hermes-vip
Group=hermes-vip
ExecStart=/usr/local/bin/hermes-vipd
Restart=always
RestartSec=5

# macOS 踩坑记录: NoNewPrivileges=yes 会阻断 daemon 的 sudo 提权
# 安全边界在 VIP 审批卡 + stamp 验证，不在 systemd 沙箱
RuntimeDirectory=hermes-vip
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
